# ali-ddns/ddns_update.service
[Unit]
Description=Aliyun ddns
After=network.target

[Service]
ExecStart=/usr/local/bin/ddns_update.py

[Install]
WantedBy=multi-user.target
